# assets/config/supervisor/6-express-proxy.conf
[program:Express-Proxy]
command=node server.js
directory=/server
priority=6
autorestart=true
stderr_logfile=/var/log/express-proxy.err.log
stdout_logfile=/var/log/express-proxy.out.log